# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=1MB
logfile_backups=1
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:nginx]
command=/usr/sbin/nginx -c /tmp/nginx.conf -g "daemon off;"
priority=10
autostart=true
autorestart=true
startsecs=2
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worldmonitor-api]
command=node /app/local-api-server.mjs
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=2
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
